# scripts/crontab
* * * * * /usr/local/bin/python /app/main.py help >> /app/mainjob.log 2>&1
